# ai-necklace.service
[Unit]
Description=AI Necklace Voice Assistant
After=network-online.target sound.target
Wants=network-online.target sound.target
# Wi-Fi接続が確立されるまで待機
Requires=network-online.target

[Service]
Type=simple
User=hotaka
Group=hotaka
WorkingDirectory=/home/hotaka/ai-necklace
Environment=PATH=/home/hotaka/ai-necklace/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/hotaka/ai-necklace/venv/bin/python3 /home/hotaka/ai-necklace/ai_necklace.py
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10
StandardOutput=journal
StandardError=journal

# オーディオデバイスへのアクセス許可
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
